# services/rtorrent/run
#!/bin/bash

dtach -n /.dtach_rtorrent /usr/bin/rtorrent -n -o import=/config/rtorrent/rtorrent.rc
inotifywait -e delete /.dtach_rtorrent
